# wildfly.conf
[program:wildfly]
command=/opt/jboss/wildfly-14.0.1.Final/bin/standalone.sh -c standalone.xml -b 0.0.0.0 -bmanagement 0.0.0.0
stdout_events_enabled=true
stderr_events_enabled=true